# distro-setup/kube-scheduler.service
[Service]
ExecStart=
ExecStart=/usr/local/bin/kube-scheduler --config=/etc/kubernetes/kube-scheduler.yaml


